# deploy/trading-dashboard.service
[Unit]
Description=Trading Bot Web Dashboard
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
Environment=PATH=/home/ubuntu/trading-bot/venv/bin
ExecStart=/home/ubuntu/trading-bot/venv/bin/python run_dashboard.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
